# log/var/lib/dpkg/info/droplet-agent.conffiles
/etc/init/droplet-agent.conf
/etc/systemd/system/droplet-agent.service
