# log/var/lib/dpkg/info/droplet-agent.list
/.
/etc
/etc/systemd
/etc/systemd/system
/etc/systemd/system/droplet-agent.service
/etc/init
/etc/init/droplet-agent.conf
/usr
/usr/share
/usr/share/doc
/usr/share/doc/droplet-agent
/usr/share/doc/droplet-agent/changelog.Debian.gz
/usr/share/doc/droplet-agent/changelog.gz
/opt
/opt/digitalocean
/opt/digitalocean/droplet-agent
/opt/digitalocean/droplet-agent/scripts
/opt/digitalocean/droplet-agent/scripts/update.sh
/opt/digitalocean/bin
/opt/digitalocean/bin/droplet-agent
